# brightness.service
[Unit]
Description=brightness

[Service]
Type=oneshot
ExecStart=/opt/bin/brightness.sh

[Install]
WantedBy=multi-user.target
